# docker/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
autostart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

#[program:user-account-consumer]
#command=/srv/api/bin/console messenger:consume amqp_user_account --bus=messenger.bus.user_account --no-debug --no-interaction --limit=10 --time-limit=3600
#autostart=true
#autorestart=true
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0
